# infra/systemd/fence-api-v2.service
[Unit]
# Multi-user web-app FastAPI backend (Sprint 4-era).
# Runs alongside the legacy `fence-api.service` (port 8503) — different
# port so both can coexist during cutover. After the new stack is
# verified for one prod day, retire `fence-api.service` and the
# Streamlit `fence-fast.service` (see deployment_notes.md).
Description=Fence Detector v2 — multi-user FastAPI backend
After=network-online.target
Wants=network-online.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/leo_streamlit
EnvironmentFile=/etc/fence-api-v2/env
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/ubuntu/leo_streamlit/venv/bin/uvicorn api_server:app \
  --host 127.0.0.1 \
  --port 8513 \
  --workers 1 \
  --proxy-headers \
  --forwarded-allow-ips=127.0.0.1
Restart=always
RestartSec=3
LimitNOFILE=65535
NoNewPrivileges=true

# Match the existing fence-api.service ceiling so MuPDF on dense pages
# doesn't OOM the box. Tune up if `journalctl -u fence-api-v2 | grep MemoryMax`
# shows kills.
MemoryMax=12G
MemoryHigh=10G

# Crashes from MuPDF or fitz can take down the worker; just respawn.
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
